# roles/nginx-letsencrypt/files/acmetiny.timer
[Unit]
Description=Renew our Let's Encrypt certificate every month

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
