# roles/nginx-letsencrypt/files/acmetiny.timer
[Unit]
Description=Renew our Let's encrypt certificate every month

[Timer]
OnCalendar=monthly
Persistent=true
